# entry.crontab
*/2 *	*	*	*	/home/eacousineau/devel/git/ramfs-test/sync.sh

