# entry.crontab
*/10 *	*	*	*	$HOME/devel/git/ramdisk-setup/sync.sh
